# configs/.config/autostart/startup.sh.desktop
[Desktop Entry]
Exec=/home/chxei/.config/conky/startup.sh
Icon=
Name=startup.sh
Path=
Terminal=False
Type=Application
